# packaging/prerm
#!/bin/sh
set -e
systemctl stop sccs-lcd.service || true
systemctl disable sccs-lcd.service || true
exit 0
